# lib/deprec/templates/god/upstart.conf
description	"God process monitor"

start on runlevel [2345]
stop on runlevel [06]

respawn
expect fork

exec /usr/local/rvm/wrappers/ree-1.8.7-2011.03/god -c /etc/god.god --log-level info -l /var/log/god.log --no-syslog